# systemd/additional/atlasvm-network-rollback.service
[Unit]
Description=AtlasVM Management Network Rollback

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/atlasvm-network-rollback
